# iparser/crontab
* * * * * /home/ivan/Project/Iparser/venv/bin/python /home/ivan/Project/Iparser/iparser/manage.py instaparser_parse_from_username >> instaparser_parse_from_username.log 2>&1
